# deploy/quadlet/syslog-hose.container
[Unit]
Description=syslog-hose — syslog traffic generator

[Container]
# Built locally with podman (see this dir's README); no registry is used.
Image=localhost/syslog-hose:latest
ContainerName=syslog-hose
Network=yard.network
PublishPort=8080:8080
Volume=hose-data.volume:/data
Environment=HOSE_ADDR=:8080
Environment=HOSE_DATA=/data
Environment=HOSE_SUGGESTED_DEST=syslog-valve:514
Environment=YARD_LINK_HOSE=8080
Environment=YARD_LINK_VALVE=8081
Environment=YARD_LINK_BUCKET=8082
# Guard this UI with the bucket's user accounts (docs/AUTH.md); remove to run open.
Environment=YARD_AUTH_URL=http://syslog-bucket:8080

[Service]
Restart=always

[Install]
WantedBy=default.target
